# install/rc.local
.........
........
#sudo nano /etc/rc.local

#add next line to  /etc/rc.local right before "exit 0"


/home/pi/teddy-vision/start.sh &


exit 0